# systemd/telegram.service
[Unit]
Description=MiiVPN - Telegram
After=network.target

[Service]
Type=simple
User=komkov.vv
WorkingDirectory=/home/komkov.vv/mii_api
ExecStart=/home/komkov.vv/mii_api/tg_handler
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target